# pkgs/workbench/cldemo-wbench-sflow-base/debian/DEBIAN/prerm
#!/bin/bash
update-rc.d -f sflowovsd remove

exit 0
