# pkgs/workbench/cldemo-wbench-sflow-base/debian/DEBIAN/prerm
#!/bin/bash
update-rc.d -f sflowovsd remove
update-rc.d -f sflow-rt remove

exit 0
